# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=/var/run/supervisord.pid

; ── Apache ────────────────────────────────────────────────────────────────────
[program:apache]
command=/usr/sbin/apache2ctl -D FOREGROUND
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/braintrust/apache.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/var/log/braintrust/apache-error.log
stderr_logfile_maxbytes=5MB

; ── WebSocket server (chat notifications) ────────────────────────────────────
; Ports 8081 (client WS) + 8082 (internal PHP notify HTTP)
[program:ws-server]
command=node /var/www/html/collabchat/ws_server.js
directory=/var/www/html/collabchat
autostart=true
autorestart=true
priority=20
user=www-data
stdout_logfile=/var/log/braintrust/ws-server.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/var/log/braintrust/ws-server-error.log
stderr_logfile_maxbytes=5MB

; ── Terminal WebSocket server (xterm.js real PTY) ─────────────────────────────
; Port 8083
[program:terminal-ws]
command=node /var/www/html/collabchat/terminal_ws.js
directory=/var/www/html/collabchat
autostart=true
autorestart=true
priority=20
; Needs root for PTY creation
user=root
environment=BT_PROJECTS_ROOT="/var/www/html/collabchat/projects/"
stdout_logfile=/var/log/braintrust/terminal-ws.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/var/log/braintrust/terminal-ws-error.log
stderr_logfile_maxbytes=5MB

; ── AI Manager (persistent Claude/Gemini CLI session orchestrator) ─────────────
[program:ai-manager]
command=node /var/www/html/collabchat/ai_manager.js
directory=/var/www/html/collabchat
autostart=true
autorestart=true
priority=20
user=www-data
stdout_logfile=/var/log/braintrust/ai-manager.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/braintrust/ai-manager-error.log
stderr_logfile_maxbytes=5MB
